# crontab
# Run every 3 hours
0 */3 * * * /scripts/download_patreon.sh >> /downloads/cron.log 2>&1
# Empty line at end required for cron
