# crontab
# Run every 6 hours
0 */6 * * * /scripts/download_patreon.sh >> /downloads/cron.log 2>&1
# Empty line at end required for cron
